# layout/DEBIAN/postinst
#!/bin/sh
launchctl unload /System/Library/LaunchDaemons/com.bigboss.sbsettingsd.plist
launchctl load   /System/Library/LaunchDaemons/com.bigboss.sbsettingsd.plist
